# server_settings/django_supervisor.conf
[program:uwsgi]
command = uwsgi --ini /opt/django/server_settings/uwsgi.ini --touch-reload /opt/django/server_settings/reload

[program:nginx]
command = /usr/sbin/nginx
